# ui/go-client/raspberrypi/pi-bluetooth/wifi_bt.service
[Unit]
Description=Bluetooth WiFi Provision
After=bluetooth.service

[Service]
ExecStart=/usr/bin/python3 /home/assmin/wifi_receiver.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
